# main/mysql/supervisord.conf
[program:mysql]
command=mysqld
autorestart=true
user=mysql
